# deploy/ubuntu/goal-execution.service.in
[Unit]
Description=goal_execution REST API (uvicorn)
After=network.target skstudio.service

[Service]
Type=simple
User=ge
Group=ge
WorkingDirectory=__APP_ROOT__
EnvironmentFile=/etc/goal-execution/goal-execution.env
ExecStart=__APP_ROOT__/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port __PORT__
Restart=on-failure
RestartSec=3
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
